# overlay/etc/systemd/system/pusher-queue.path
[Unit]
Description=Manages the pusher queue on disk

[Path]
DirectoryNotEmpty=/run/sensar/pusher-queue
Unit=pusher.service
MakeDirectory=true
