# overlay/etc/systemd/system/sensar.timer
[Unit]
Description=the timer that controls execution of sensar

[Timer]
OnCalendar=*-*-* *:*:00
Unit=sniffer.service
