# overlay/etc/systemd/system/sensar.timer
[Unit]
Description=the timer that controls execution of sensar

[Timer]
OnActiveSec=5min
Unit=sniffer.service
