# overlay/etc/systemd/system/sensar.target
[Unit]
Wants=pusher-queue.path
Wants=pusher.service
Wants=sniffer.service
