# overlay/etc/systemd/system/sensar.target
[Unit]
After=network.target
Wants=pusher-queue.path
Wants=pusher.service
Wants=sensar.timer
Wants=sniffer.service
