# overlay/etc/systemd/system/sniffer.service
[Unit]
Description=Collects data for sensar.target

[Service]
Type=oneshot
ExecStart=-/usr/bin/sniff
