# overlay/etc/systemd/system/sniffer.service
[Unit]
Description=Collects data for sensar.target

[Service]
ExecStart=/usr/bin/sniff
